# .docker/dev/php/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
logfile_backups=3
logfile_maxbytes=50MB
user=root

[unix_http_server]
file=/var/run/supervisor.sock;

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock;

[program:php]
command=docker-php-entrypoint php-fpm
numprocs=1
autostart=true
autorestart=true

[program:cron]
command=/usr/sbin/crond -f -l 2
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/cron-error.log
stdout_logfile=/var/log/supervisor/cron.log
stdout_logfile_backups=3
stderr_logfile_backups=3